# ev3dev-jessie/root/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.

# ------------------------------------------------------------------------------
# If there's a file called ev3dev.rc.local on the FAT32 partition, run it

if [ -e /media/mmc_p1/ev3dev.rc.local ]; then
    echo "Executing /media/mmc_p1/ev3dev.rc.local"
    . /media/mmc_p1/ev3dev.rc.local
fi

# ------------------------------------------------------------------------------
# Turn on the cursor for the virtual consoles (the cursor is turned off in the
# kernel so that it does not interfere with the boot splash).
export TERM=linux
echo 1 > /sys/module/vt/parameters/global_cursor_default
ACTIVE_CONSOLES="/dev/tty[1-7]"
for tty in $ACTIVE_CONSOLES; do
    setterm -cursor on > $tty
done


# ------------------------------------------------------------------------------
# Give us something pretty to look at until we get brickdm working
chvt 1

exit 0

