# default/hooks/firstboot
cat > $ROOTDIR/etc/init.d/firstboot << __END__
#!/bin/sh -e
### BEGIN INIT INFO
# Provides: firstboot
# Required-Start: \$all
# Required-Stop:
# Default-Start: S
# Default-Stop:
### END INIT INFO

update-rc.d -f firstboot remove

# generating ssh host key to be unique
for f in rsa dsa ecdsa ed25519; do

    if [ -f /etc/ssh/ssh_host_\${f}_key ]; then
        rm -f /etc/ssh/ssh_host_\${f}_key
        echo "generating ssh \${f} key..."
        ssh-keygen -q -t \${f} -f /etc/ssh/ssh_host_\${f}_key -C "" -N ""
    fi
done
__END__

chmod +x $ROOTDIR/etc/init.d/firstboot
$CHROOTQEMUCMD $ROOTDIR update-rc.d firstboot start 99 S
